# stratipi/etc/cron.d/logrotate
@daily		root	/usr/local/sbin/logrotate	/usr/local/etc/logrotate.conf

