# stratipi/etc/cron.d/resize
@reboot		root	/bin/sh /sbin/zresize.sh

