# stratipi/etc/cron.d/save_time
@daily		root	/usr/bin/touch /var/db/last_time

