# Files/etc/cron.d/save_time
@daily		root	/usr/bin/touch /etc/last_time
